# systemd-config/kbot-chatbot.service
[Unit]
Description=Directline Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/usr/local/bin/integra-kbot/chatbot
ExecStart=/usr/local/bin/integra-kbot/chatbot/IndianBank_ChatBOT
Restart=always
RestartSec=3s

[Install]
WantedBy=default.target
